# restic/systemd/daily-backup@.timer
[Unit]
Description=timer for daily backup of %i

[Timer]
OnCalendar=daily
Persistent=true
Unit=run-backup@%i.service

[Install]
WantedBy=default.target
