# restic/systemd/run-backup@.service
[Unit]
Description=backup of %i

[Service]
Type=oneshot
ExecStart=/root/bin/backups/call_restic %i.env backup
